# chezmoi/private_Library/private_LaunchAgents/com.koopa.color-mode-sync.plist.tmpl
{{- if eq .chezmoi.os "darwin" -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.koopa.color-mode-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>{{ .chezmoi.homeDir }}/.local/share/koopa/bin/koopa configure user color-mode</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>{{ .chezmoi.homeDir }}/Library/Preferences/.GlobalPreferences.plist</string>
  </array>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/koopa-color-mode-sync.log</string>
  <key>StandardErrorPath</key>
  <string>{{ .chezmoi.homeDir }}/Library/Logs/koopa-color-mode-sync.log</string>
</dict>
</plist>
{{- end }}
